# ansible/templates/project_backend.service.j2
[Unit]
Description=Project {{ group_code }} Flask Backend Service (Automated via Ansible)
After=postgresql.service
Requires=postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory={{ backend_dir }}
Environment="DB_HOST=localhost"
Environment="DB_NAME={{ db_name }}"
Environment="DB_USER={{ db_user }}"
Environment="DB_PASS={{ db_password }}"
ExecStart={{ backend_dir }}/venv/bin/python app.py
Restart=always

[Install]
WantedBy=multi-user.target